# watch_xchange.service
[Unit]
Description=xchange restart daemon for Status NG
After = xchange.service

[Service]
ExecStart=python /usr/local/bin/watch_xchange.py
Restart=always
Type=simple

[Install]
WantedBy=multi-user.target
